# rc.d/aged
#!/bin/sh
#
# PROVIDE: aged
# REQUIRE: DAEMON
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable aged:
#
# aged_enable="YES"
# aged_flags=""           # Additional flags (e.g., "-v" for verbose)
# aged_socket="/var/run/aged.sock"
# aged_database="/var/db/aged"
#

. /etc/rc.subr

name="aged"
rcvar="aged_enable"

load_rc_config $name

: ${aged_enable:="NO"}
: ${aged_socket:="/var/run/aged.sock"}
: ${aged_database:="/var/db/aged"}

pidfile="/var/run/${name}.pid"
command="%%PREFIX%%/bin/aged"
command_args="-f -s ${aged_socket} -d ${aged_database} ${aged_flags}"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

aged_start()
{
    if aged_status_quiet; then
        echo "${name} is already running."
        return 1
    fi

    echo "Starting ${name}."
    # Use daemon(8) to:
    #   -p: write PID file
    #   -r: restart on exit (supervised mode)
    #   -S: syslog output
    # aged runs with -f (foreground) so daemon(8) can manage it
    /usr/sbin/daemon -p ${pidfile} -r -S -T ${name} ${command} ${command_args}
}

aged_stop()
{
    if [ -f ${pidfile} ]; then
        pid=$(cat ${pidfile})
        echo "Stopping ${name}."
        kill -TERM ${pid} 2>/dev/null
        wait_for_pids ${pid}
        rm -f ${pidfile}
    else
        echo "${name} is not running."
        return 1
    fi
}

aged_status()
{
    if [ -f ${pidfile} ]; then
        pid=$(cat ${pidfile})
        if kill -0 ${pid} 2>/dev/null; then
            echo "${name} is running as pid ${pid}."
            return 0
        else
            echo "${name} is not running (stale pidfile)."
            rm -f ${pidfile}
            return 1
        fi
    else
        echo "${name} is not running."
        return 1
    fi
}

aged_status_quiet()
{
    if [ -f ${pidfile} ]; then
        pid=$(cat ${pidfile})
        kill -0 ${pid} 2>/dev/null
        return $?
    fi
    return 1
}

run_rc_command "$1"
